# eww-clean.timer
[Unit]
Description=CyberClean Auto Timer — every 6h

[Timer]
OnBootSec=5min
OnUnitActiveSec=6h
Persistent=true

[Install]
WantedBy=timers.target
